# tools/ttm.upstart
description "ttm server"

start on runlevel [2345]
stop on runlevel [!2345]

pre-start script
mkdir -p /var/log/ttm
chown ttm:root /var/log/ttm/
mkdir -p /var/run/ttm
chown ttm:root /var/run/ttm/

mkdir -p /var/lock/ttm
chown ttm:root /var/lock/ttm/
end script


exec start-stop-daemon --start --chuid ttm --chdir /opt/stack/ttm -m --pidfile /var/run/ttm.pid --name ttm --exec /opt/stack/ttm/tools/with_venv.sh -- python /opt/stack/ttm/ttm.py > /var/log/ttm/ttm.log 2>&1
