# tools/ttm.upstart
description "ttm server"

start on runlevel [2345]
stop on runlevel [!2345]

pre-start script
mkdir -p /var/run/ttm
chown root /var/run/ttm/

mkdir -p /var/lock/ttm
chown root /var/lock/ttm/
end script

exec start-stop-daemon --start --chdir /opt/stack/ttm --exec /opt/stack/ttm/tools/with_venv.sh -- python /opt/stack/ttm/ttm.py
